# linux/hdmi_sound.rules
SUBSYSTEM=="drm", RUN+="/usr/local/bin/hdmi_sound_toggle"
